# run-latest-build-systemd@.service
[Service]
Type=simple
ExecStart=/bin/run-latest-build.sh %i

[Install]
WantedBy=multi-user.target
